# scripts/systemd/market-data.service
[Unit]
Description=Market Data Service
After=network.target

[Service]
WorkingDirectory=/mnt/market_data
Environment=RUST_LOG=info
ExecStart=/mnt/market_data/target/release/market-data
Restart=always
RestartSec=5
TimeoutStopSec=300
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
